# twitch-recorder.service
[Unit]
Description=Twitch Stream Recorder Daemon
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/twitch-ui-dowload
ExecStart=/usr/bin/python3 /path/to/twitch-ui-dowload/daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
